# bin/etc/systemd/eard.service
[Unit]
Description=EARD - Energy Aware Runtime node daemon
ConditionPathExists=/var/ear
#Wants=req1.service req2.service
#After=req1.service req2.service


[Service]
Type=simple
Environment="EAR_ETC=/root/ear/bin/etc"
ExecStart=/root/ear/bin/sbin/eard 1
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
KillMode=process
Restart=on-failure
#RestartSec=1s
LimitMEMLOCK=infinity
LimitSTACK=infinity

[Install]
WantedBy=multi-user.target

